# DEBIAN/postrm
#!/bin/sh
update-rc.d usb-time-card remove
